# src/autoshutdown.service
[Unit]
Description=Autoshutdown

[Service]
ExecStart=/usr/bin/autoshutdown.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target